# systemd/sage-router.service
[Unit]
Description=Sage Router - Intent-based AI model routing
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
EnvironmentFile=-%h/.config/sage-router/sage-router.env
ExecStart=/usr/bin/env bash -lc 'cd "$${SAGE_ROUTER_HOME:?Set SAGE_ROUTER_HOME in ~/.config/sage-router/sage-router.env}" && if [ -n "$${SAGE_ROUTER_PATH_PREFIX:-}" ]; then export PATH="$${SAGE_ROUTER_PATH_PREFIX}:$$PATH"; fi; exec "$${SAGE_ROUTER_PYTHON:-python3}" router.py --port "$${SAGE_ROUTER_PORT:-8788}" $${SAGE_ROUTER_ARGS:-}'
Restart=always
RestartSec=10

[Install]
WantedBy=default.target
